# deploy/launchd/ai.hermes.signet.mcp.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.hermes.signet.mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/SIGNET/.venv/bin/signet</string>
    <string>deployment</string>
    <string>serve-mcp</string>
    <string>--config</string>
    <string>/ABSOLUTE/PATH/TO/SIGNET-DATA/config/disabled.json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/SIGNET</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>Umask</key>
  <integer>63</integer>
  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/SIGNET-DATA/logs/mcp.log</string>
  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/SIGNET-DATA/logs/mcp-error.log</string>
</dict>
</plist>
